# redis/files/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf

