# mongodb/files/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:mongod]
command=/usr/bin/mongod --config /etc/mongodb.conf

